# systemd/icarus-prune.service
[Unit]
Description=Icarus weekly wiki prune
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/jeon/icarus
EnvironmentFile=/home/jeon/icarus/.env
Environment=PATH=/home/jeon/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=CLAUDE_BIN=/home/jeon/.local/bin/claude
ExecStart=/usr/bin/node --dns-result-order=ipv4first /home/jeon/icarus/node_modules/.bin/tsx /home/jeon/icarus/scripts/weekly-prune.ts
StandardOutput=append:/home/jeon/icarus/logs/prune.log
StandardError=append:/home/jeon/icarus/logs/prune.log
