# systemd/icarus.service
[Unit]
Description=Icarus (Telegram life agent)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/jeon/icarus
EnvironmentFile=/home/jeon/icarus/.env
Environment=PATH=/home/jeon/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=CLAUDE_BIN=/home/jeon/.local/bin/claude
ExecStart=/usr/bin/node --dns-result-order=ipv4first /home/jeon/icarus/dist/index.js
Restart=on-failure
RestartSec=5
StandardOutput=append:/home/jeon/icarus/logs/stdout.log
StandardError=append:/home/jeon/icarus/logs/stderr.log

[Install]
WantedBy=default.target
