# systemd/icarus-prune.service
[Unit]
Description=Icarus weekly wiki prune
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/jeon/Projects/icarus
EnvironmentFile=/home/jeon/Projects/icarus/.env
Environment=PATH=/home/jeon/.nvm/versions/node/v22.22.2/bin:/home/jeon/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=CLAUDE_BIN=/home/jeon/.local/bin/claude
ExecStart=/home/jeon/Projects/icarus/node_modules/.bin/tsx /home/jeon/Projects/icarus/scripts/weekly-prune.ts
StandardOutput=append:/home/jeon/Projects/icarus/logs/prune.log
StandardError=append:/home/jeon/Projects/icarus/logs/prune.log
